# overseer.service
[Unit]
Description=Overseer - Nginx Log Watcher for Nydus
After=network.target mysql.service mariadb.service

[Service]
Type=simple
User=overseer
Group=overseer
SupplementaryGroups=www-data
WorkingDirectory=/opt/overseer
ExecStartPre=/bin/mkdir -p /states
ExecStartPre=/bin/chown overseer:overseer /states
ExecStartPre=/bin/chmod 755 /states
ExecStart=/opt/overseer/venv/bin/python main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
MemoryMax=500M
MemorySwapMax=0
ReadWritePaths=/states

[Install]
WantedBy=multi-user.target